# scripts/core_services/cat-ink-syncer.service
[Unit]
Description=LeopardCat AgentOS Session Syncer (Operation Cat-Ink)
After=network.target

[Service]
EnvironmentFile=-%h/agentmanager/.env
WorkingDirectory=%h/agentmanager
ExecStart=/usr/bin/python3 %h/agentmanager/scripts/core_services/session_syncer.py
Restart=always
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
